# gtfs.conf
description "gtfs server"
author "Aoyama, Shotaro"
start on runlevel [2345]
stop on runlevel [016]
exec sudo -i -u skybrain bash -c 'cd /home/skybrain/gtfs && bin/start'
respawn
